# pastebin.service
[Unit]
Description=Pastebin Node.js Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/pastebin
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target

